# emulated-hue-dev/rootfs/etc/services.d/emulated-hue/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Emulated Hue
# Runs Emulated Hue
# ==============================================================================

export ESPHOME_IS_HASSIO=true

if bashio::config.has_value 'data'; then
    export DATA_DIR=$(bashio::config 'data')
fi

if bashio::config.has_value 'url'; then
    export HASS_URL=$(bashio::config 'url')
fi

if bashio::config.has_value 'token'; then
    export HASS_TOKEN=$(bashio::config 'token')
fi

if bashio::config.has_value 'verbose'; then
    export VERBOSE=$(bashio::config 'verbose')
fi

bashio::log.info "Starting Emulated Hue..."
exec python3 -m /app/emulated_hue
